# cloudformation/etcd/etc/systemd/system/etcd-init.service
[Unit]
Description=etd init step

[Service]
Type=oneshot
ExecStart=/usr/local/bin/etcd-init.sh

[Install]
WantedBy=multi-user.target
